# deploy/space/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

# The TTS router (internal, :8080). Runs from source via Bun.
[program:router]
command=bun /app/apps/router/src/index.ts
environment=PORT="8080",NODE_ENV="production"
autorestart=true
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# The Next.js web app on the port HF exposes (:7860).
[program:web]
command=bun /app/standalone/apps/web/server.js
directory=/app/standalone
environment=PORT="7860",HOSTNAME="0.0.0.0",NODE_ENV="production",ROUTER_URL="http://127.0.0.1:8080",AUDIO_LOG_DIR="/audio",SQLITE_PATH="/data/tts_arena.db"
autorestart=true
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
